# evaluation/supervisor/imo-2025-problem1-eval.conf
[program:imo-2025-problem1-eval]
environment=PROC_NAME="%(program_name)s"
command=/opt/supervisor-scripts/imo-2025-problem1-eval.sh
autostart=false
autorestart=false
startsecs=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
stopsignal=TERM
